# bin/init/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /clonecademy/django/uwsgi.ini --logto /tmp/uwsgi.log

[program:nginx-app]
command = /usr/sbin/nginx
